# cookbooks-sources/cube/templates/default/upstart.conf.erb
description "Cube <%= @component %>"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

env PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

chdir <%= node['cube']['root'] %>

exec sudo node bin/<%= @component %>.js
